# netatmo/api/systemd/netatmo-api.service
# /etc/systemd/system/netatmo.service

[Unit]
Description=get netatmo data from cloud api

[Service]
User=simon
Type=oneshot
ExecStart=/usr/bin/python3m /home/simon/netatmo/api/netatmo.py

[Install]
WantedBy=multi-user.target
